# gear_version/src/service/autolock.service
[Unit]
Description=autolock
After=pigpiod.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/takumi/autolock/gear_version
ExecStart=/home/takumi/autolock/gear_version/autolock
ExecStop=/bin/kill -WINCH ${MAINPID}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

